# com.onpoint.antigravity-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onpoint.antigravity-cdp</string>

    <!-- Launch Antigravity with CDP enabled on port 9222 -->
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Antigravity.app/Contents/MacOS/Electron</string>
        <string>--remote-debugging-port=9222</string>
        <string>--remote-allow-origins=*</string>
    </array>

    <!-- Start automatically at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Only run when user is logged in (GUI session) -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>BRIDGE_DIR_PLACEHOLDER/antigravity-cdp.log</string>

    <key>StandardErrorPath</key>
    <string>BRIDGE_DIR_PLACEHOLDER/antigravity-cdp-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>HOME_PLACEHOLDER</string>
    </dict>
</dict>
</plist>
